# cron.d/root
0 */1 * * * /data/adb/modules/android-fstrim/auto_trim.sh
