# cron.d/root
0 */2 * * * /data/adb/modules/android-fstrim/auto_trim.sh
